# deploy/cleaners-hub-autodeploy.timer
[Unit]
Description=Poll for cleaners-hub deploys every minute

[Timer]
# Fire 30s after the timer is enabled (whether at boot or via install-vps.sh
# bootstrap), then every 60s after the previous run. OnBootSec alone would
# only fire shortly after boot; with OnActiveSec we also catch fresh
# enable-after-install events. AccuracySec is loose so neighboring timers
# can batch in the same wake-up.
OnActiveSec=30s
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=cleaners-hub-autodeploy.service

[Install]
WantedBy=timers.target
